# reading-helper.service
[Unit]
Description=Reading Helper Server
After=network.target

[Service]
Type=simple
User=wangchengan
WorkingDirectory=/Users/wangchengan/Documents/reading-helper
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=CONFIG_DIR=/Users/wangchengan/Documents/reading-helper/config
Environment=USER_DATA_ROOT=/Users/wangchengan/Documents/reading-helper/data/users
ExecStart=/opt/homebrew/bin/node /Users/wangchengan/Documents/reading-helper/server/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
