# reading-helper.service
[Unit]
Description=Reading Helper Server
After=network.target

[Service]
Type=simple
User=<your-user>
WorkingDirectory=/path/to/reading-helper
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=CONFIG_DIR=/path/to/reading-helper/config
Environment=USER_DATA_ROOT=/path/to/reading-helper/data/users
ExecStart=/usr/bin/node /path/to/reading-helper/server/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
